# system_files/etc/systemd/system/tailscaled.service.d/flags.conf
[Service]
Environment=FLAGS=--netfilter-mode=nodivert --no-logs-no-support
