# dotfiles/dot_config/systemd/user/mpvpaper-dp3.service
[Unit]
Description=mpvpaper wallpaper on DP-3
PartOf=graphical-session.target
After=graphical-session.target

[Service]
ExecStart=mpvpaper -o "no-audio --loop-file=100 --reset-on-next-file=loop-file --loop-playlist shuffle --profile=fast" DP-3 %h/videos/wallpapers/landscape/
Restart=on-failure
TimeoutStopSec=5

[Install]
WantedBy=graphical-session.target
