# dotfiles/dot_config/systemd/user/wallpaper-cycle.timer
[Unit]
Description=Periodically rotate per-monitor wallpapers

[Timer]
OnBootSec=15s
OnUnitActiveSec=30m

[Install]
WantedBy=timers.target
