# dotfiles/dot_config/systemd/user/wallpaper-theme.timer
[Unit]
Description=Periodically check the wallpaper video for theme changes

[Timer]
OnBootSec=15s
OnUnitActiveSec=30s

[Install]
WantedBy=timers.target
