# dotfiles/dot_config/systemd/user/wallpaper-cycle.service
[Unit]
Description=Rotate per-monitor wallpapers and regenerate the matugen theme
PartOf=graphical-session.target
After=graphical-session.target awww-daemon.service

[Service]
Type=oneshot
ExecStart=%h/.local/bin/wallpaper-cycle
