# dotfiles/dot_config/systemd/user/wallpaper-theme.service
[Unit]
Description=Regenerate matugen theme from the current mpvpaper wallpaper frame
PartOf=mpvpaper-dp2.service

[Service]
Type=oneshot
ExecStart=%h/.local/bin/wallpaper-theme
